# .Services/gustave.target
[Unit]
Description=Starts all gustave services
Wants=gustave.service gustave-jobs.service gustave-scheduler.service